# supervisor.ini
[program:twitter]
command=/usr/bin/python3 /home/wodim/akari_endlosung/twitter_bot.py
directory=/home/wodim/akari_endlosung
user=wodim
autorestart=true
autostart=true

[program:telegram]
command=/usr/bin/python3 /home/wodim/akari_endlosung/telegram_bot.py
directory=/home/wodim/akari_endlosung
user=wodim
autorestart=true
autostart=true

[group:akari]
programs=twitter,telegram
priority=999
