# supervisor.ini
[program:twitter]
command=/usr/bin/python3.5 /home/wodim/akari_endlosung/twitter_bot.py
directory=/home/wodim/akari_endlosung
user=wodim
autorestart=true
autostart=true
environment=LC_ALL="en_US.UTF-8"

[program:telegram]
command=/usr/bin/python3.5 /home/wodim/akari_endlosung/telegram_bot.py
directory=/home/wodim/akari_endlosung
user=wodim
autorestart=true
autostart=true
environment=LC_ALL="en_US.UTF-8"

[group:akari]
programs=twitter,telegram
priority=999
